# ej/serial-app/build/deploy-runtime-window/mxt-runtime.service
[Unit]
Description=MXT Runtime Window HTTP (Python)
After=network.target

[Service]
Type=simple
Environment=HTTP_PORT=__HTTP_PORT__
Environment=WEB_ROOT=__WEB_ROOT__
ExecStart=/usr/bin/python3 /opt/mxt-runtime/serve_runtime.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
